# vzfailcnt.cron
*/1 * * * * root /usr/sbin/vzfailcnt conf
