# vzfailcnt.cron
# Put this file to /etc/cron.d/
*/1 * * * * root /usr/sbin/vzfailcnt
